# server/systemd/sparqplug-updater.service
[Unit]
Description=SparQ stack updater (pull repos, rebuild containers)
Wants=docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/sparqplug/server
Environment=REPO_ROOT=/opt/sparqplug/server
ExecStart=/usr/bin/env bash -lc "/usr/local/bin/sparqplug-update.sh"
User=root
Group=root

[Install]
WantedBy=multi-user.target
